# systemroot/etc/inittab
::sysinit:/etc/init.d/rcS
::respawn:/sbin/getty 115200 tty1
::ctrlaltdel:/sbin/reboot
::respawn:/usr/sbin/dropbear -F -E